# user/xdg-document-portal.service.d/override.conf
[Service]
UMask = 0077
IPAddressDeny = any
DevicePolicy=closed
## ProtectClock = yes BREAKS!
IPAddressDeny = any
