# system/systemd-homed.service.d/override.conf
[Service]
ProtectKernelModules = yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
ProtectControlGroups = yes
ProtectKernelLogs = yes
ProtectSystem = strict
