# system/logrotate.service.d/override.conf
[Service]
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
UMask = 0077
IPAddressDeny = any
SystemCallArchitectures=native
SystemCallFilter=@system-service
